# conf/site.conf
[program:site]
command=/home/bd/venv/bin/uwsgi --ini uwsgi.ini
directory=/home/bd/site
stopsignal=QUIT
stdout_logfile=/home/bd/tmp/site.stdout.log
stderr_logfile=/home/bd/tmp/site.stderr.log